# deployment-scripts/supervisor_mdpsapp.conf
[program:mdps]
command = /var/lib/jenkins/workspace/gunicorn_start.bash                  ; Command to start app
user = jenkins                                                ; User to run as
stdout_logfile = /var/lib/jenkins/workspace/logs/gunicorn_supervisor.log   ; Where to write log messages
redirect_stderr = true                                       ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8              ; Set UTF-8 as default encoding